# pi/deploy/dog-door.service
[Unit]
Description=Dog door web + Arduino serial bridge
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/arduino-dog-door-automation/pi
Environment=PYTHONPATH=/home/pi/arduino-dog-door-automation/pi
EnvironmentFile=/home/pi/arduino-dog-door-automation/pi/.env
ExecStart=/home/pi/arduino-dog-door-automation/pi/.venv/bin/uvicorn dog_door_pi.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
